# raspberry_pi/scripts/spotless-kiosk.service
[Unit]
Description=Project Spotless - Kiosk Web Server
After=network-online.target mosquitto.service
Wants=network-online.target mosquitto.service

[Service]
Type=simple
User=spotless
Group=spotless
WorkingDirectory=/home/spotless/IOT_Spotless/raspberry_pi
EnvironmentFile=-/home/spotless/IOT_Spotless/raspberry_pi/.env
ExecStart=/home/spotless/IOT_Spotless/raspberry_pi/venv/bin/python /home/spotless/IOT_Spotless/raspberry_pi/main.py --kiosk
Restart=always
RestartSec=5
StartLimitIntervalSec=120
StartLimitBurst=10
StandardOutput=journal
StandardError=journal

# Shutdown behaviour. The main process installs SIGTERM handlers that
# turn off relays, close MQTT, then os._exit(0). Werkzeug otherwise
# never returns on its own, so without this systemd waited the full
# TimeoutStopSec on every reboot. 15s is plenty for the cleanup; if
# anything genuinely hangs we'd rather KILL than block the reboot.
KillMode=mixed
KillSignal=SIGTERM
SendSIGKILL=yes
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
